# fomalhaut/deploy/fomalhaut.ini
[program:fomalhaut]
command=python -m fomalhaut.runserver --port=651%(process_num)d
numprocs=4
process_name=%(program_name)s_651%(process_num)d
directory=/home/python
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisord/fomalhaut.log
loglevel=info
user=fomalhaut
